# crontab
# Syncs container with git repo
*/30 * * * * root bash git-sync.sh
